# etc/systemd/system/margin-fetch.service
[Unit]
Description=JPX 信用取引残高データ日次取得 (stock-dialy)
After=network.target

[Service]
Type=oneshot
User=naoki
Group=naoki
WorkingDirectory=/var/www/django/stock-dialy
Environment="DJANGO_SETTINGS_MODULE=config.settings"

# 過去40日分をチェック（未取得の日があれば取得）
ExecStart=/var/www/django/stock-dialy/venv/bin/python manage.py fetch_margin_data --days 40

# ログはjournaldに出力 (journalctl -u margin-fetch で確認可能)
StandardOutput=journal
StandardError=journal
